# nfs_project/mount_files/home-skogsparlan-photos.automount
[Unit]
Description=Foton och filmklipp på NAS
ConditionPathExists=/home/skogsparlan/photos

[Automount]
Where=/home/skogsparlan/photos
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
