# nfs_project/mount_files/home-skogsparlan-data.automount
[Unit]
Description=Generella datafiler på NAS
ConditionPathExists=/home/skogsparlan/data

[Automount]
Where=/home/skogsparlan/data
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
