# nfs_project/mount_files/home-skogsparlan-emp.automount
[Unit]
Description=Datafiler på NAS
ConditionPathExists=/home/skogsparlan/emp

[Automount]
Where=/home/skogsparlan/emp
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
