# nfs_project/mount_files/home-skogsparlan-backup.automount
[Unit]
Description=Backup datafiler på NAS
ConditionPathExists=/home/skogsparlan/backup

[Automount]
Where=/home/skogsparlan/backup
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
